# ns-microservices/ns-users/dist/debian/DEBIAN/postinst
#!/bin/bash
sudo ln -sf /etc/nginx/sites-available/ns_gateway.conf /etc/nginx/sites-enabled/ns_gateway.conf
chown -R root:root /usr/bin/{NAME}
systemctl daemon-reload
# systemctl enable {NAME}_{PORT_1}.service
# systemctl start {NAME}_{PORT_1}.service
# systemctl enable {NAME}_{PORT_2}.service
# systemctl start {NAME}_{PORT_2}.service
# systemctl enable {NAME}_{PORT_3}.service
# systemctl start {NAME}_{PORT_3}.service
# echo ""
# echo "The service {NAME} has been installed."
# echo "The following instances are running:"
# echo "    nsusers_{PORT_1} at address 0.0.0.0:{PORT_1}"
# echo "    nsusers_{PORT_2} at address 0.0.0.0:{PORT_2}"
# echo "    nsusers_{PORT_3} at address 0.0.0.0:{PORT_3}"
# echo ""
# echo "Please restart nginx to make the service API reachable."
# echo ""
